# systemd/zram-setup.service
[Unit]
Description=ZRAM compressed swap setup
Before=swap.target
Wants=swap.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/zram-setup.sh start
ExecStop=/usr/local/bin/zram-setup.sh stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
